# LaunchAgent/com.acuity.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Reference only. `acuity install` generates the authoritative agent at
  ~/Library/LaunchAgents/com.acuity.agent.plist via AgentManager.buildPlist,
  resolving the real Acuity.app binary path on this machine.

  Key requirements (do not regress):
   - Launch the binary INSIDE Acuity.app (bundle identity → WindowServer access).
     A bare CLI binary (e.g. /usr/local/bin/acuity) exits EX_CONFIG (78) under launchd.
   - KeepAlive=true so the menubar relaunches on quit/crash (cannot silently vanish).
   - LimitLoadToSessionType=Aqua; load via `launchctl bootstrap gui/<uid>`.
   - Log to ~/Library/Logs (user-owned), never /tmp (root-owned footgun — see GOTCHAS.md).
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.acuity.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOU/Applications/Acuity.app/Contents/MacOS/acuity</string>
        <string>start</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOU/Library/Logs/acuity.stderr.log</string>
</dict>
</plist>
